# init.d/minecraft-branch
#!/sbin/runscript
# Minecraft: Branch mini-world
# vim:ft=gentoo-init-d

depend() {
	need minecraft-any
}

start() {
	: ;
}

stop() {
	: ;
}
