# init.d/minecraft-bungeecord
#!/sbin/runscript
# Minecraft: BungeeCord multiplexer proxy
# vim:ft=gentoo-init-d

MC_UID=`id -u minecraft`

depend() {
	need minecraft-branch minecraft-laboratory minecraft-caterpillars minecraft-butterflies
}

start() {
	su -c '/usr/bin/tmux -L Minecraft new-window -n BungeeCord -c /srv/minecraft/BungeeCord /usr/bin/java -Xmx512M -jar /opt/minecraft/BungeeCord.jar' minecraft
	# TODO: watch logs or output for: Done (N.NNNs)! For help, type "help"
}

stop() {
	/usr/bin/tmux -S /tmp/tmux-$MC_UID/Minecraft send-keys -t Minecraft:BungeeCord e n d Enter
	# TODO: watch logs or output for successful stop
	#   timeout: kill tmux window if it's still around
}
