# init.d/minecraft-laboratory
#!/sbin/runscript
# Minecraft: Laboratory creative-mode flat world
# vim:ft=gentoo-init-d

depend() {
	need minecraft-any
}

start() {
	: ;
}

stop() {
	: ;
}
