# init.d/minecraft-all
#!/sbin/runscript
# Minecraft: BungeeCord multiplexer proxy
# vim:ft=gentoo-init-d

depend() {
	need minecraft-bungeecord
}

start() {
	: ;
}

stop() {
	: ;
}
