# init.d/minecraft-branch
#!/sbin/runscript
# Minecraft: Branch mini-world

depend() {
	need minecraft-any
}

start() {
	: ;
}

stop() {
	: ;
}
