# init.d/minecraft-any
#!/sbin/runscript
# Does nothing. Here to be a required dependency of all the Minecraft worlds.
# So that I can stop them all at once, just by stopping this one.

depend() {
	need net
	use dns
}

start() {
	: ;
}

stop() {
	: ;
}
