# init.d/minecraft-caterpillars
#!/sbin/runscript
# Minecraft: Caterpillars survival-mode (easy) world

depend() {
	need minecraft-any
}

start() {
	: ;
}

stop() {
	: ;
}
