# init.d/minecraft-laboratory
#!/sbin/runscript
# Minecraft: Laboratory creative-mode flat world

depend() {
	need minecraft-any
}

start() {
	: ;
}

stop() {
	: ;
}
