# init.d/minecraft-butterflies
#!/sbin/runscript
# Minecraft: Butterflies survival-mode (hard) world

depend() {
	need minecraft-any
}

start() {
	: ;
}

stop() {
	: ;
}
