# init.d/minecraft-all
#!/sbin/runscript
# Minecraft: BungeeCord multiplexer proxy

depend() {
	need minecraft-bungeecord
}

start() {
	: ;
}

stop() {
	: ;
}
