# init.d/minecraft-any
#!/sbin/runscript
# Does nothing. Here to be a required dependency of all the Minecraft worlds.
# So that I can stop them all at once, just by stopping this one.
# vim:ft=gentoo-init-d

MC_UID=`id -u minecraft`

depend() {
	need net
	use dns
}

start() {
	checkpath -m 666 /dev/ptmx
	su -c '/usr/bin/tmux -L Minecraft new-session -d -s Minecraft -c /srv/minecraft' minecraft
}

stop() {
	/usr/bin/tmux -S /tmp/tmux-$MC_UID/Minecraft kill-session -t Minecraft
	: ;
}
