# init.d/minecraft-caterpillars
#!/sbin/runscript
# Minecraft: Caterpillars survival-mode (easy) world
# vim:ft=gentoo-init-d

command='/usr/bin/tmux'
command_args='-L Minecraft new-window -n Caterpillars /usr/bin/java -Xms512M -Xmx2048M -jar /opt/minecraft/spigot-1.8.8.jar'
start_stop_daemon_args='--user minecraft:minecraft --chdir /srv/minecraft/Caterpillars'

depend() {
	need minecraft-any
}
